# supervisord.conf
; =============================================================================
; supervisord.conf — runs Haggai's desktop stack inside the container.
; supervisord itself runs as root (so it can setuid each program to `user`);
; every actual program runs as the unprivileged `user`. Start order is by
; `priority`:  Xvfb (10) -> pulseaudio (15) -> XFCE+dbus (20) -> rustdesk (30).
; The XFCE and RustDesk scripts each wait for the X server before proceeding and
; fail loud if it never appears.
; =============================================================================
[supervisord]
nodaemon=true
user=root
loglevel=info
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

[program:xvfb]
command=/usr/local/bin/start-xvfb.sh
user=user
priority=10
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
environment=HOME="/home/user",USER="user",XDG_RUNTIME_DIR="/run/user/1000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:pulseaudio]
command=/usr/local/bin/start-pulseaudio.sh
user=user
priority=15
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
environment=HOME="/home/user",USER="user",XDG_RUNTIME_DIR="/run/user/1000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:xfce]
command=/usr/local/bin/start-xfce.sh
user=user
priority=20
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
environment=HOME="/home/user",USER="user",DISPLAY=":99",XDG_RUNTIME_DIR="/run/user/1000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:rustdesk]
command=/usr/local/bin/start-rustdesk.sh
user=user
priority=30
autorestart=true
startretries=1000000
stopasgroup=true
killasgroup=true
environment=HOME="/home/user",USER="user",DISPLAY=":99",XDG_RUNTIME_DIR="/run/user/1000"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
